# etc/cron.d/yabt
# YABT main cron job
*/10 * * * * root /usr/local/bin/yabt
